# dashboard/dashboard.service
[Unit]
Description=Dashboard Fede — FastAPI system & OpenRouter monitor
After=network-online.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fede-infra/dashboard
ExecStart=/usr/bin/python3 /opt/fede-infra/dashboard/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# OpenRouter keys via environment
# Set CLAW_OPENROUTER_API_KEY manually or use an EnvironmentFile
EnvironmentFile=/root/.openclaw-ticia/env.sh
Environment=AUTO_TOPUP=on
Environment=AGENT_TELEGRAM=@s_ticia_bot

[Install]
WantedBy=multi-user.target